# contrib/build-openwrt-kamikazeipk/wifidog/files/wifidog.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=65
EXTRA_COMMANDS="status"
EXTRA_HELP="        status Print the status of the service"


start() {
	/usr/bin/router-plus-portal-init start
}

stop() {
	/usr/bin/router-plus-portal-init stop
}

status() {
	/usr/bin/router-plus-portal-init status
}